# agents/rocky/ashd-agent.service
[Unit]
Description=ASHD Monitoring Agent
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/ashd-agent
ExecStart=/usr/bin/python3 /opt/ashd-agent/ashd_agent.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
